# root/etc/services.d/plex-web/run
#!/usr/bin/with-contenv bash

# Run flask server
exec \
  s6-setuidgid abc flask run