# root/etc/services.d/plex-web/run
# Run flask server
exec \
	s6-setuidgid abc flask run